# systemd/headphone-hum-mask.service
[Unit]
Description=Quiet headphone hum to mask idle DAC noise
After=sound.target pipewire-pulse.service pulseaudio.service
Wants=sound.target

[Service]
Type=simple
ExecStart=%h/.local/bin/headphone-hum-mask.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
